# docker/supervisord.conf
; =============================================================================
; Oxidized Bio - Supervisor Configuration
; =============================================================================
; Manages multiple processes within the container:
;   - oxidized-bio: Main Rust application (HTTP API + RFC endpoints)
;   - sshd: OpenSSH server for RFC remote shell access
;   - log-cleanup: Periodic log rotation
;
; Supervisor provides:
;   - Process monitoring and automatic restart on failure
;   - Graceful shutdown handling
;   - Centralized logging
;   - Process group management
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info

; Socket for supervisorctl
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

; RPC interface for supervisorctl
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; CLI control interface
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; =============================================================================
; Oxidized Bio Main Application
; =============================================================================
; The main Rust application handling HTTP requests and RFC calls
; Environment variables are inherited from container - no explicit mapping needed
[program:oxidized-bio]
command=/app/oxidized-bio --server
directory=/app
user=oxidized
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
killasgroup=true
stopasgroup=true

; Logging configuration
stdout_logfile=/app/logs/oxidized-bio.stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/oxidized-bio.stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5

; Redirect stderr to stdout for unified logging (optional)
; redirect_stderr=true

; Priority determines startup order (lower = starts first)
priority=100

; =============================================================================
; OpenSSH Server
; =============================================================================
; Provides SSH access for RFC remote shell commands
; Can be disabled via ENABLE_SSH=false environment variable
[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
startsecs=2
startretries=3
stopwaitsecs=10

; Logging
stdout_logfile=/app/logs/sshd.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/sshd.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

; Start SSH before main app
priority=50

; =============================================================================
; Log Cleanup Job
; =============================================================================
; Runs periodically to clean up old log files
[program:log-cleanup]
command=/bin/bash -c "while true; do find /app/logs -name '*.log.*' -mtime +7 -delete 2>/dev/null; find /var/log/supervisor -name '*.log.*' -mtime +7 -delete 2>/dev/null; sleep 86400; done"
autostart=true
autorestart=true
startsecs=0
user=root

; Don't log the cleanup job itself
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; Lowest priority (starts last)
priority=999

; =============================================================================
; Process Groups
; =============================================================================
; Groups allow controlling multiple processes together
[group:oxidized]
programs=oxidized-bio,sshd,log-cleanup
priority=100

; =============================================================================
; Event Listeners (Optional)
; =============================================================================
; Uncomment to enable email notifications on process crashes
; [eventlistener:crashmail]
; command=/usr/local/bin/crashmail -a -m admin@example.com
; events=PROCESS_STATE_EXITED
; buffer_size=100
